# root/etc/s6-overlay/s6-rc.d/init-nomachine/run
#!/bin/sh
dbus-launch && /etc/NX/nxserver --startup
